# stockmarket.service
[Unit]
Description=StockMarket Platform (Go + HTMX)
After=network.target

[Service]
Type=simple
User=exedev
WorkingDirectory=/home/exedev/stockmarket
ExecStart=/home/exedev/stockmarket/bin/server
Environment=PORT=5000
Environment=DATABASE_PATH=/home/exedev/stockmarket/stockmarket.db
Environment=ENVIRONMENT=production
Environment=ENCRYPTION_KEY=+wTmXVVXvlEbz1J8cxFeTc2D53Z+xxS0AJpsbIKrYcM=
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
